# deploy/supervisor.conf
[program:pushmatrix]
directory=/www/wwwroot/PushMatrix/backend
command=/www/wwwroot/PushMatrix/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/var/log/pushmatrix.err.log
stdout_logfile=/var/log/pushmatrix.out.log
user=www
environment=PYTHONPATH="/www/wwwroot/PushMatrix/backend"

[program:pushmatrix-celery]
directory=/www/wwwroot/PushMatrix/backend
command=/www/wwwroot/PushMatrix/backend/venv/bin/celery -A celery_app worker --loglevel=info
autostart=true
autorestart=true
startsecs=5
stderr_logfile=/var/log/pushmatrix-celery.err.log
stdout_logfile=/var/log/pushmatrix-celery.out.log
user=www
